# 17ce
#!/bin/sh /etc/rc.common
# Copyright (C) 2017  www.17ce.com
START=99
CDN_BASE="http://git.oschina.net/mjyhj/K2_17ce/raw/master/"
WORK_DIR="/tmp/17ce/"
init_conf()
{
#	config_load '17ceclient'
#	config_get username 17ceclient username

	if [ ! -f "/tmp/lib/libstdc++.so.6" ]; then
	wget -O /tmp/lib/libstdc++.so.6   $CDN_BASE/lib/libstdc++.so.6.0.21
	wget -O /tmp/17ce/17ce_v3      $CDN_BASE/bin/17ce_v3
	wget -O /tmp/17ce/conf.json    $CDN_BASE/bin/conf.json
	fi

	if [ -f "/tmp/17ce/UUID" ]; then cd /tmp/17ce; rm UUID 17ce_v3.log user; fi;
}

start()
{
  sleep 5
	init_conf
  ln -sf /lib/libc.so.0  /tmp/lib/libc.so
  ln -sf /us/lib/libpolassl.so  /tmp/lib/libmbedtls.so.9
	$WORK_DIR/17ce_v3 -u USER_NAME
	echo "17ce Client has started."
}

stop()
{
	killall -9 17ce_v3
	sleep 1
	echo "17ce Client has stoped."
}
